# chezmoi/private_dot_config/private_aconfmgr/files/private_etc/systemd/system/private_reflector.timer.d/override.conf
[Timer]
OnCalendar=monthly
Persistent=true
